# scripts/crontab
# Crontab for LLM-GEN-REPORT Backend automation
# Format: minute hour day month weekday command

# Daily cleanup at 2:00 AM (remove files older than 7 days)
0 2 * * * cd /app && python scripts/cleanup.py >> /app/data/logs/cleanup.log 2>&1

# Daily backup at 3:00 AM (create backup of uploads)
0 3 * * * cd /app && python scripts/backup.py >> /app/data/logs/backup.log 2>&1

# Weekly storage report every Sunday at 1:00 AM
0 1 * * 0 cd /app && python scripts/cleanup.py >> /app/data/logs/storage-report.log 2>&1

# Empty line required at end of crontab
